# couchdb/service/couchdb/run
#!/bin/sh
exec /sbin/setuser couchdb /usr/bin/couchdb >> /var/log/couchdb.log 2>&1