# wsgi/service/pymc/run
#!/bin/sh
exec /sbin/setuser pymc /usr/bin/uwsgi /var/lib/pymc/uwsgi.ini \
        >> /var/lib/pymc/pymc.io.log 2>&1
